# runtime/static-web.conf
[program:lighttpd]
directory = /home/admin/runtime
command = /home/admin/runtime/lighttpd/bin/lighttpd -f /home/admin/runtime/lighttpd/conf/lighttpd.conf -m /home/admin/runtime/lighttpd/lib  -D
process_name = lighttpd
stdout_logfile = /var/log/bae/lighttpd.stdout.log
stderr_logfile = /var/log/bae/lighttpd.stderr.log
autorestart = true

